# private_dot_config/profile.d/40_path.sh
# User bin
pathmunge "$HOME"/.local/bin
