# private_dot_config/profile.d/90_application_defaults.sh
## Application defaults
export EDITOR="nano"
export BROWSER="firefox"

# Set less to skip the pager if one page and to allow control codes (e.g. colors)
export LESS="-RF"
# Set gnu time format to minute-level ISO ('+%Y-%m-%d %H:%M')
export TIME_STYLE="long-iso"
